# service/reefwatch.service
[Unit]
Description=ReefWatch Security Daemon (HIDS)
After=network.target
Documentation=https://github.com/yasnaak/reefwatch

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m reefwatch.daemon --config %h/.openclaw/workspace/skills/reefwatch/reefwatch_config.yaml
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.openclaw %t
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
LockPersonality=true
RestrictSUIDSGID=true
PrivateDevices=true
SystemCallFilter=@system-service

[Install]
WantedBy=default.target
